# files/supervisord/conf.d/foreman.conf
[program:foreman]
user=foreman
directory=/usr/share/foreman
command=/usr/bin/bundle exec bin/rails server -b 0.0.0.0 -p 443 -e production